# crontab
# Execute o script Python diariamente às 21h
0 21 * * * root cd /app && git pull && /usr/local/bin/python /app/rel_finan_google_old.py >> /var/log/cron.log 2>&1
